# dashboard/dashboard-ws.service
#
# dashboard service test
# https://github.com/ess-dmsc/integration-test/dashboard
#

[Unit]
Description=Dashboard for ECDC server monitoring
After=network.target

[Service]
User=ecdcdashboard
WorkingDirectory=/var/opt/dm_group/dashboard
ExecStart=/usr/bin/python3 -m SimpleHTTPServer 8765
